# scripts/com.sentrial.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sentrial.menubar</string>

    <!-- Launch via the Sentrial.app bundle so macOS attributes microphone +
         Input Monitoring permission requests to Sentrial (whose Info.plist
         carries NSMicrophoneUsageDescription and NSInputMonitoringUsageDescription)
         rather than to Python.app (whose signed Info.plist lacks both and
         therefore causes silent TCC denial for the bare python3 binary). -->
    <key>ProgramArguments</key>
    <array>
        <string>__SENTRIAL_HOME__/Sentrial.app/Contents/MacOS/Sentrial</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SENTRIAL_HOME__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__SENTRIAL_LOGS__/menubar.out.log</string>

    <key>StandardErrorPath</key>
    <string>__SENTRIAL_LOGS__/menubar.err.log</string>

    <!-- menubar must run as a LaunchAgent, not LaunchDaemon -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
